# systemd/spark-dashboard.service
[Unit]
Description=SparkView Web Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/lesha/spark-dashboard
Environment=PYTHONPATH=/home/lesha/sparkview
ExecStart=/home/lesha/spark-dashboard/.venv/bin/uvicorn spark_dashboard.app:app --host 0.0.0.0 --port 8500 --no-access-log
Restart=always
RestartSec=5

[Install]
WantedBy=default.target