# systemd/spark-dashboard.service
[Unit]
Description=SparkView Web Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/spark-dashboard
Environment=PYTHONPATH=%h/sparkview
ExecStart=%h/spark-dashboard/.venv/bin/uvicorn spark_dashboard.app:app --host 0.0.0.0 --port 8500 --no-access-log
Restart=always
RestartSec=5

[Install]
WantedBy=default.target